# configs/supervisor/conf.d/gunicorn.conf
[program:gunicorn]
command = /home/pi/www/flask_rest_dht/env/bin/gunicorn -w 1 main:app
directory = /home/pi/www/flask_rest_dht
user = pi
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/gunicorn.log
stderr_logfile = /var/log/supervisor/gunicorn_err.log
